# debian/vanilla-systemd-units.postinst
#!/bin/bash
umount -l /usr
units="\x2esystem-usr.mount \x2esystem-usr-local.mount"

for unit in $units; do
  echo "Enabling unit file: $unit..."
  systemctl enable $unit
done

systemctl daemon-reload

umount -l /usr  # workaround: dpkg still needs to be able to access /usr to complete the transaction
echo "Reboot to complete installation."
true